# systemd/removepii-streamlit.service
[Unit]
Description=RemovePII Streamlit frontend
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=vpc
Group=vpc
WorkingDirectory=/home/vpc/projects/removepii
Environment=PATH=/home/vpc/projects/removepii/.venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/home/vpc/projects/removepii/.env
ExecStart=/home/vpc/projects/removepii/.venv/bin/streamlit run frontend.py \
    --server.port=8501 \
    --server.address=127.0.0.1 \
    --server.headless=true \
    --browser.gatherUsageStats=false
Restart=on-failure
RestartSec=5

# Harden a bit — this thing sees paste input; assume hostile data.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/vpc/projects/removepii

[Install]
WantedBy=multi-user.target
